# upstart_script/cheese.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn
setuid group
chdir /home/group/Dailydos
exec ./script.sh