# upstart_script/cheese.conf
start on runlevel [2345]
stop on runleve [!2345]

respawn
exec /home/group/Dailydos/script.sh